# scripts/Deployment/database_access.service.sample
[Unit]
Description=database_access.py service file
After=network.target

[Service]
User=
WorkingDirectory=
ExecStart=
Restart=always

[Install]
WantedBy=multi-user.target
